# deployment/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:pibasho-backend]
command=pibasho-backend
autorestart=unexpected
exitcodes=0,2
stopsignal=TERM